# systemd/user/battery-monitor.service
[Unit]
Description=Monitors the battery Power and State.

[Service]
Type=simple
ExecStart=%h/.config/hypr/battery
Restart=on-failure
RestartSec=5
