# systemd/user/ntfy_sub.service
[Unit]
Description=Ntfy Subscribe Notifications

[Service]
Type=simple
ExecStart=%h/.config/hypr/ntfy_sub
Restart=always
RestartSec=5
