# systemd/user/power-mode@.service
[Unit]
Description=Power Mode (%i)

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/battery %i
