# systemd/user/power-mode@.service
[Unit]
Description=Power Mode (%i)

[Service]
Type=oneshot
ExecStart=/home/tanveer/.config/hypr/battery %i
