# systemd/user/theme.service
[Unit]
Description=Set GTK Theme

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/themes
